# docker/s6/services/.s6-svscan/finish
#!/bin/sh

echo "Container shutting down"